# systemd/weekly-discord-bot.service
[Unit]
Description=Weekly Transaction Pipeline Discord Bot
After=network.target

[Service]
Type=simple
# User akan diisi otomatis oleh install_service.sh
User=<USER>
# WorkingDirectory akan diisi otomatis oleh install_service.sh
WorkingDirectory=<PROJECT_DIR>
ExecStart=/usr/bin/env node index.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=weekly-discord-bot

# Loads environment variables from the .env file
EnvironmentFile=<PROJECT_DIR>/.env

[Install]
WantedBy=multi-user.target
